# agent/com.clrc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clrc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/clrc</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/clrc.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/clrc.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/opt/homebrew/sbin</string>
        <key>TERM</key>
        <string>xterm-256color</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>

        <!-- Generate once: uuidgen | tr A-Z a-z -->
        <key>AGENT_ID</key>
        <string>YOUR_UUID</string>

        <!-- Display name shown in the iOS app -->
        <key>AGENT_NAME</key>
        <string>My Mac</string>

        <!-- Must match AGENT_SECRET in relay/.env -->
        <key>AGENT_SECRET</key>
        <string>YOUR_AGENT_SECRET</string>

        <!-- ws://LAN_IP:8080  or  wss://your-relay.up.railway.app -->
        <key>RELAY_URL</key>
        <string>ws://192.168.1.10:8080</string>

        <key>DEFAULT_COMMAND</key>
        <string>claude</string>
    </dict>
</dict>
</plist>
